# etc/profile.d/jupyterlab.sh
if [[ -z "$LAB_PATH" ]]; then
    export LAB_PATH="###LAB_PATH###"
fi